# udev-tv.rules
SUBSYSTEM=="usb" \
, ACTION=="add" \
, ATTR{idProduct}=="006d" \
, ATTR{idVendor}=="256c" \
, RUN+="/usr/local/bin/tv start"
