# systemd/keepassxc.service
[Unit]
After=xdg-desktop-autostart.target

[Service]
Type=simple
ExecStart=/usr/bin/keepassxc

[Install]
WantedBy=graphical-session.target
